# OSDevGrp.OSIntranet.WebApi.supervisord.conf
[supervisord]
nodaemon=true

[program:webapi]
command=/usr/bin/dotnet /app/OSDevGrp.OSIntranet.WebApi.dll

[program:ssh]
command=/usr/sbin/sshd -D